# .config/autostart/autocutsel.desktop
[Desktop Entry]
Name=autocutsel
Exec=autocutsel
Terminal=false
Type=Application
StartupNotify=false
