# .config/autostart/edit-server.desktop
[Desktop Entry]
Name=Edit Server
Exec=/home/kill/bin/edit-server
Terminal=false
Type=Application
StartupNotify=false
