# .config/autostart/mpd.desktop
[Desktop Entry]
Name=MPD Server
Exec=mpd ~/.mpdconf
Terminal=false
Type=Application
StartupNotify=false
