# .config/autostart/xcompmgr.desktop
[Desktop Entry]
Name=Xcompmgr
Exec=xcompmgr
Terminal=false
Type=Application
StartupNotify=false
